# apps/install/chimera.service
[Unit]
Description=Chimera Local AI Node
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'cd "${HOME}/.local/share/chimera/qvac" && bash start-auto.sh'
Restart=on-failure
RestartSec=5
Environment=PORT=3002
Environment=LOG_LEVEL=info

[Install]
WantedBy=default.target
